# ansible/roles/proxmox/networking/files/rebolt.service
[Unit]
Description=Check FRR topology

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rebolt.sh